# ntopng/build/supervisord.conf
[supervisord]
nodaemon=true

[program:redis-server]
command=/usr/bin/redis-server /etc/redis/redis.conf --daemonize no 
user=redis

[program:nprobe]
command=/usr/local/bin/nprobe /etc/nprobe/nprobe-none.conf

[program:ntopng]
command=/usr/local/bin/ntopng %(ENV_NTOPARGS)s
