# nfsen/build/supervisord.conf
[supervisord]
nodaemon = true

[program:syslogd]
command = /sbin/syslogd -n -l '%(ENV_LOGLEVEL)s' -O /proc/1/fd/1
priority = 1

[program:nfsen]
command = /data/nfsen/bin/nfsen start
redirect_stderr = true
stdout_logfile = syslog

[program:php-fpm]
command = /usr/local/sbin/php-fpm --nodaemonize --fpm-config /conf/php-fpm.conf
redirect_stderr = true
stdout_logfile = syslog

[program:nginx]
command = /usr/sbin/nginx -c /conf/nginx.conf -g 'daemon off;'
redirect_stderr = true
stdout_logfile = syslog
